# raid-notifier.service
[Unit]
Description=Raid Notifier
Requires=dbus.service

[Service]
Type=simple
ExecStart=/usr/bin/raid-notifier

[Install]
WantedBy=default.target
